# systemd/erp-saas-autodeploy.service
[Unit]
Description=ERP SaaS auto deploy from GitHub
After=network-online.target docker.service
Wants=network-online.target docker.service

[Service]
Type=oneshot
ExecStart=/usr/local/bin/erp-saas-autodeploy.sh